# deploy/systemd/infx-llm.service
[Unit]
Description=infx worker
After=network-online.target
Wants=network-online.target

[Service]
User=infx
Group=infx
EnvironmentFile=-/etc/infx/worker.env
ExecStart=/usr/bin/infx-llm $WORKER_FLAGS
Restart=on-failure
RestartSec=3s
AmbientCapabilities=
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
CapabilityBoundingSet=
StateDirectory=infx
RuntimeDirectory=infx
# Hardening knobs (adjust as needed)
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
