# debian/nfrx-asr.service
[Unit]
Description=nfrx asr worker
After=network-online.target
Wants=network-online.target

[Service]
User=nfrx
Group=nfrx
EnvironmentFile=-/etc/nfrx/asr.env
ExecStart=/usr/bin/nfrx-asr $ASR_FLAGS
Restart=on-failure
RestartSec=3s
AmbientCapabilities=
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
CapabilityBoundingSet=
StateDirectory=nfrx
RuntimeDirectory=nfrx
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
